# scripts/hawk-api.service
[Unit]
Description=Hawk News Service translation API
After=network.target postgresql.service redis.service

[Service]
Type=simple
User=jamditis
WorkingDirectory=/home/jamditis/projects/hawk-translation-api
EnvironmentFile=/home/jamditis/projects/hawk-translation-api/.env
ExecStart=/home/jamditis/projects/hawk-translation-api/venv/bin/uvicorn api.main:app --host 127.0.0.1 --port 8090
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
