# scripts/hawk-worker.service
[Unit]
Description=Hawk News Service Celery worker
After=network.target redis.service

[Service]
Type=simple
User=joe
WorkingDirectory=/home/joe/projects/hawk-translation-api
EnvironmentFile=/home/joe/projects/hawk-translation-api/.env
ExecStart=/home/joe/projects/hawk-translation-api/venv/bin/celery -A workers.celery_app worker --loglevel=info --concurrency=2
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
